# root/etc/services.d/dillinger/run
#!/usr/bin/with-contenv bash

cd /app/dillinger || exit

exec \
	s6-setuidgid abc npm start
