# root/etc/services.d/dillinger/run
#!/usr/bin/with-contenv bash

cd /app/dillinger || exit

NODE_ENV=production exec \
	s6-setuidgid abc node app
