# root/etc/services.d/dillinger/run
#!/usr/bin/with-contenv bash

if [[ "$(uname -m)" == "armv7l" ]]; then
    sleep infinity
fi

cd /app/dillinger || exit

NODE_ENV=production exec \
	s6-setuidgid abc node app
